# infrastructure/cron/nightly_bitget_checks.cron
# Nightly Bitget stop-loss validation
# Установить через: crontab infrastructure/cron/nightly_bitget_checks.cron
# Требуется активное виртуальное окружение .venv и доступ к Bitget ключам в env.

0 2 * * * cd /Users/zhuchyok/Documents/GITHUB/atra/atra && /bin/bash scripts/run_nightly_bitget_checks.sh >> logs/test_results.log 2>&1

